# docker/crontab
*/5 * * * * root /usr/local/bin/node /src/api/cli.js push-task sites.scanAll >/dev/null 2>&1

* 2 * * * root /usr/local/bin/node /src/api/cli.js push-task pages.scan >/dev/null 2>&1